# sandbox/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[program:xvfb]
priority=10
command=/usr/bin/Xvfb :1 -screen 0 1280x900x24 -ac +extension GLX +render -noreset
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:chromium]
priority=20
command=/bin/bash /opt/sandbox/start-chromium.sh
environment=DISPLAY=":1",PLAYWRIGHT_BROWSERS_PATH="/home/sandbox/.cache-pw"
autorestart=true
startretries=5
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:socat]
priority=30
command=/usr/bin/socat TCP-LISTEN:9222,fork,reuseaddr TCP:127.0.0.1:8222
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:x11vnc]
priority=40
command=/usr/bin/x11vnc -display :1 -nopw -forever -shared -rfbport 5900 -xkb -ncache 10
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:websockify]
priority=45
command=/usr/bin/websockify 0.0.0.0:5901 localhost:5900
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:sandbox-api]
priority=50
command=/usr/bin/python3 -m uvicorn api.main:app --host 0.0.0.0 --port 8080
directory=/opt/sandbox
environment=PLAYWRIGHT_BROWSERS_PATH="/home/sandbox/.cache-pw"
autorestart=true
startretries=10
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
